# database_container/crontab
*/15	*	*	*	*	/bin/sh /vuls/go-cve-dictionary.sh
